# wayland/services/keyring/run
#!/bin/execlineb -P
gnome-keyring-daemon --foreground --components=pkcs11,secrets
