# wayland/services/ssh-agent/run
#!/bin/execlineb -P
importas -i dir XDG_RUNTIME_DIR
foreground { mkdir -p ${dir}/ssh }
ssh-agent -a ${dir}/ssh/agent -t 7200 -D
