# wayland/services/xps/run
#!/bin/execlineb -P
/usr/libexec/xdg-permission-store -r
