# wayland/services/pulse/run
#!/bin/execlineb -P
pulseaudio --daemonize=no --exit-idle-time=-1
