# wayland/services/wireplumber/run
#!/bin/execlineb -P
tryexec
{ wireplumber }
pipewire-media-session
