# wayland/services/xdp/run
#!/bin/execlineb -P
/usr/libexec/xdg-desktop-portal -r
