# wayland/services/spotifyd/run
#!/bin/execlineb -P
spotifyd --no-daemon
