# wayland/services/swayidle/run
#!/bin/execlineb -P
swayidle -w
timeout 300  "swaylock"
before-sleep "swaylock"
