# void/void.d/etc-x86/sv/kmscon/run
#!/bin/execlineb -P
kmscon --vt=tty3 --switch-vt
