# wayland/services/xdocp/run
#!/bin/execlineb -P
/usr/libexec/xdg-document-portal -r
