# wayland/services/gb/run
#!/bin/execlineb -P
hittpd -p 2024 build/
# add build/ symlink to glowing bear directory
